# etc/crontabs/root
*/1 * * * * /bin/sh /root/check.sh > /dev/null 2>&1
